# packaging/systemd/iot-lwm2m-client.service
[Unit]
Description=IoT LwM2M client (binds to /run/iot DM server URI)
Documentation=file:///etc/iot/lwm2m-client.env
# Needs the data-store socket for live config; ordering ensures
# DsConfig.connect() doesn't fall through to defaults at boot.
After=network-online.target iot-ds.service
Wants=network-online.target iot-ds.service
# StartLimit* live in [Unit] for modern systemd (250+).
StartLimitBurst=5
StartLimitIntervalSec=60s

[Service]
Type=simple
# All knobs in the env file so operators edit /etc/iot/lwm2m-client.env,
# not the unit. Missing file → unit fails fast with "EnvironmentFile not
# found" (the `-` prefix below makes it optional; we want it required so
# typos surface, hence no prefix).
EnvironmentFile=/etc/iot/lwm2m-client.env
ExecStart=/usr/local/bin/lwm2m $LWM2M_ARGS

# PSK provisioning (apps/docs/tdd-psk-provisioning.md): the client runs
# as the STATIC `engineer` account (NOT DynamicUser) so it has a stable
# primary gid the data-store read_acl/write_acl (gid:engineer) can match
# — that is what lets only this client read the write-only PSK keys.
# `iot` (supplementary) grants access to ds-server's 0660 socket. The
# account is provisioned by packaging/systemd/iot-sysusers.conf.
User=engineer
Group=engineer
SupplementaryGroups=iot
# No RuntimeDirectory here — iot-ds owns /run/iot (0755, traversable);
# the client only connects to the socket. After=iot-ds.service orders it.

# A PSK change in dev-mode makes the client self-exit; systemd relaunches
# it with fresh credentials. Restart=always so a clean exit also restarts.
Restart=always
RestartSec=5s

StandardOutput=journal
StandardError=journal

NoNewPrivileges=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictRealtime=yes

[Install]
WantedBy=multi-user.target
