# packaging/systemd/iot-sensord.service
[Unit]
Description=IoT sensor producer (mangOH Yellow BMI160/BME680/light over I2C → ds iot.sensor.*)
Documentation=file:///etc/iot/sensord.env
Documentation=file:///usr/local/share/doc/iot/tdd-mangoh-yellow-sensors.md
# Needs ds-server up (publishes iot.sensor.* on connect).
After=network-online.target iot-ds.service
Wants=iot-ds.service
StartLimitBurst=5
StartLimitIntervalSec=60s

[Service]
Type=simple
EnvironmentFile=/etc/iot/sensord.env
ExecStart=/usr/local/bin/iot-sensord $IOT_SENSORD_ARGS

# PREFERRED path: the kernel i2c-dev node (/dev/i2c-1) — no CAP_SYS_RAWIO,
# just device access (boot with dtparam=i2c_arm=on so the node exists). The
# daemon falls back to the BSC register transport via /dev/mem only if i2c-dev
# is absent, which needs CAP_SYS_RAWIO (kept here for that fallback; drop it,
# and the /dev/mem/gpiomem DeviceAllow, once i2c-dev is confirmed on the board).
DeviceAllow=/dev/i2c-1 rw
DeviceAllow=char-i2c rw
AmbientCapabilities=CAP_SYS_RAWIO
CapabilityBoundingSet=CAP_SYS_RAWIO
DeviceAllow=/dev/mem rw
DeviceAllow=/dev/gpiomem rw

# /run/iot is owned solely by tmpfiles.d (iot.conf, 2775 root:iot) — do NOT use
# RuntimeDirectory=iot: it resets /run/iot to root:root 0755 on start and removes
# it on stop, clobbering the group-iot ds socket + breaking every other daemon's
# data-store access. The daemon just reads /run/iot/data_store.sock like the rest.

Restart=on-failure
RestartSec=10s

StandardOutput=journal
StandardError=journal

NoNewPrivileges=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictRealtime=yes

[Install]
WantedBy=multi-user.target
