# packaging/systemd/iot-ds.service
[Unit]
Description=IoT data-store daemon (ds-server, EMP over unix socket)
Documentation=file:///usr/local/include/data_store/proto.hpp
# Network not required — the socket is local.
After=local-fs.target
# Crash-loop protection: 3 restarts within 30 s → fail so the operator
# notices, instead of pegging a core. Lives in [Unit] not [Service]
# in modern systemd (250+).
StartLimitBurst=3
StartLimitIntervalSec=30s

[Service]
Type=simple
ExecStart=/usr/local/bin/ds-server \
    ds-socket=/run/iot/data_store.sock \
    ds-store=/var/lib/iot/data_store.lua \
    ds-schema-dir=/etc/iot/ds-schemas \
    ds-socket-group=iot

# DynamicUser auto-allocates a system user/group named `iot-ds` for the
# duration of this service; StateDirectory + RuntimeDirectory end up
# owned by that uid/gid. Cheaper than maintaining a /etc/passwd entry.
DynamicUser=yes
# PSK provisioning (task J3): join the shared `iot` group (provisioned by
# iot-sysusers.conf) so ds-server can chgrp its 0660 socket to `iot`
# (ds-socket-group= above); the static `engineer` client is also in `iot`
# and can then connect. A process may chgrp a file it owns only to a
# group it is a member of, hence this membership.
SupplementaryGroups=iot
RuntimeDirectory=iot
RuntimeDirectoryMode=0755
StateDirectory=iot
StateDirectoryMode=0750

Restart=on-failure
RestartSec=2s

# journald handles stdout/stderr (ACE logging writes to stderr) so
# `journalctl -u iot-ds` works out of the box.
StandardOutput=journal
StandardError=journal

# Modest hardening — full lockdown is FUP material (would need cap_net
# removal, ProtectSystem=strict, etc.).
NoNewPrivileges=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictRealtime=yes

[Install]
WantedBy=multi-user.target
