# yocto/meta-iot/recipes-iot/lwm2m/files/iot-factory-reset.service
[Unit]
Description=Factory reset — wipe persisted state to schema defaults, then reboot
Documentation=https://github.com/naushada/iot
# No [Install]: activated by iot-factory-reset.path, never enabled directly.

[Service]
Type=oneshot
ExecStart=/usr/bin/iot-factory-reset
TimeoutStartSec=120
