# yocto/meta-iot/recipes-iot/lwm2m/files/iot-hostname.service
[Unit]
Description=Derive iot-<serial> hostname for mDNS device-UI discovery
Documentation=https://github.com/naushada/iot
After=local-fs.target
# Set the hostname before Avahi advertises so it announces iot-<serial>.local.
Before=avahi-daemon.service
ConditionPathExists=/usr/bin/iot-set-hostname

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/iot-set-hostname
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
