# yocto/meta-iot/recipes-iot/lwm2m/files/iot-sensord.service
[Unit]
Description=IoT sensor producer (mangOH Yellow BMI160/BME680/light over I2C → ds iot.sensor.*)
Documentation=file:///etc/iot/sensord.env
# Needs ds-server up (publishes iot.sensor.* on connect).
After=network-online.target iot-ds.service
Wants=iot-ds.service
StartLimitBurst=5
StartLimitIntervalSec=60s

[Service]
Type=simple
EnvironmentFile=/etc/iot/sensord.env
ExecStart=/usr/bin/iot-sensord $IOT_SENSORD_ARGS

# PREFERRED: kernel i2c-dev node (/dev/i2c-1) — no CAP_SYS_RAWIO, just device
# access (boot with dtparam=i2c_arm=on). Falls back to the BSC /dev/mem register
# transport (needs CAP_SYS_RAWIO) only if i2c-dev is absent. Once i2c-dev is
# confirmed on the board, drop the cap + the /dev/mem,gpiomem DeviceAllow.
DeviceAllow=/dev/i2c-1 rw
DeviceAllow=char-i2c rw
AmbientCapabilities=CAP_SYS_RAWIO
CapabilityBoundingSet=CAP_SYS_RAWIO
DeviceAllow=/dev/mem rw
DeviceAllow=/dev/gpiomem rw

RuntimeDirectory=iot
RuntimeDirectoryMode=0755

Restart=on-failure
RestartSec=10s

StandardOutput=journal
StandardError=journal

NoNewPrivileges=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictRealtime=yes

[Install]
WantedBy=multi-user.target
