# yocto/meta-iot/recipes-iot/lwm2m/files/iot-reboot.service
[Unit]
Description=Reboot the device (armed by iot-reboot.path)
Documentation=https://github.com/naushada/iot
# No [Install]: activated by iot-reboot.path, never enabled directly.

[Service]
Type=oneshot
# Run via a shell so `systemctl`/`rm` resolve via PATH (matching the proven
# iot-swupdate convention) — a unit ExecStart needs an absolute path, and
# /bin/systemctl does not exist on every image (systemctl may be /usr/bin only),
# which silently no-op'd the reboot. Remove the trigger first so the .path
# unit re-arms, then reboot.
ExecStart=/bin/sh -c 'rm -f /run/iot/reboot.request; exec systemctl reboot'
