# yocto/meta-iot/recipes-iot/lwm2m/files/iot-ddnsd.service
[Unit]
Description=IoT device-side Dynamic DNS updater (keeps a public A record on the device's current public IP → ds ddns.*)
Documentation=file:///etc/iot/ddnsd.env
Documentation=file:///usr/local/share/doc/iot/tdd-ddns.md
# Needs ds-server up (reads ddns.* config, publishes ddns.state/last.*) and
# outbound network (HTTPS to the IP-echo + the DNS provider).
After=network-online.target iot-ds.service
Wants=network-online.target iot-ds.service
StartLimitBurst=5
StartLimitIntervalSec=60s

[Service]
Type=simple
# Run in group iot so the ds peer-gid matches the write-only secret keys'
# read_acl=gid:iot (ddns.<provider>.token) — the daemon can read them, ds-cli
# (root) and Viewer clients cannot. Also grants access to the group-iot ds
# socket (/run/iot/data_store.sock, 2775 root:iot).
Group=iot
EnvironmentFile=/etc/iot/ddnsd.env
ExecStart=/usr/bin/iot-ddnsd $IOT_DDNS_ARGS

# Alternative secret source (FR-9): a credential FILE instead of the write-only
# ds key. Point ddns.token.path at it, or drop it as a systemd credential (the
# daemon reads $CREDENTIALS_DIRECTORY/ddns_secret). Left commented because
# LoadCredential fails the unit if the source file is absent — enable per site:
#   LoadCredential=ddns_secret:/etc/iot/ddns/token

# /run/iot is owned solely by tmpfiles.d (iot.conf, 2775 root:iot) — do NOT use
# RuntimeDirectory=iot: it resets /run/iot on start/stop and clobbers the
# group-iot ds socket. The daemon just reads /run/iot/data_store.sock.

Restart=on-failure
RestartSec=10s

StandardOutput=journal
StandardError=journal

NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictRealtime=yes

[Install]
WantedBy=multi-user.target
