# yocto/meta-iot/recipes-iot/lwm2m/files/iot-openvpn-client.service
[Unit]
Description=IoT OpenVPN client supervisor
Documentation=https://github.com/naushada/iot
After=network-online.target iot-ds.service
Wants=network-online.target iot-ds.service
StartLimitBurst=5
StartLimitIntervalSec=60s

[Service]
Type=simple
EnvironmentFile=/etc/iot/openvpn-client.env
ExecStart=/usr/bin/openvpn-client $OPENVPN_CLIENT_ARGS

DynamicUser=yes
# Join the `iot` group to open ds-server's 0660 group-iot socket.
SupplementaryGroups=iot
AmbientCapabilities=CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_ADMIN
DeviceAllow=/dev/net/tun rw

# No RuntimeDirectory=iot: /run/iot is owned by tmpfiles.d/iot.conf. Declaring
# it here would delete the shared dir (and ds-server's socket) on stop/restart.

Restart=on-failure
RestartSec=5s

StandardOutput=journal
StandardError=journal

NoNewPrivileges=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictRealtime=yes

[Install]
WantedBy=multi-user.target
