# yocto/meta-iot/recipes-iot/lwm2m/files/iot-ds.service
[Unit]
Description=IoT data-store daemon (ds-server, EMP over unix socket)
Documentation=https://github.com/naushada/iot
After=local-fs.target
StartLimitBurst=3
StartLimitIntervalSec=30s

[Service]
Type=simple
ExecStart=/usr/bin/ds-server \
    ds-socket=/run/iot/data_store.sock \
    ds-store=/var/lib/iot/data_store.lua \
    ds-schema-dir=/etc/iot/ds-schemas \
    ds-socket-group=iot

DynamicUser=yes
# PSK provisioning (task J3): join the shared `iot` group so ds-server can
# chgrp its 0660 socket to `iot` (ds-socket-group= above); the static
# `engineer` client is also in `iot` and can then connect.
SupplementaryGroups=iot
# /run/iot is created once by tmpfiles.d/iot.conf (2775 root:iot) — deliberately
# NOT via RuntimeDirectory=iot. systemd removes a RuntimeDirectory= when the unit
# (or any other unit that also declared the same name) stops, which would delete
# the live socket out from under running clients and break ds-server's own
# restart (pre-exec dir setup churn). The setgid tmpfiles dir lets this
# DynamicUser (in group iot) create the 0660 group-iot socket.
StateDirectory=iot
StateDirectoryMode=0750

Restart=on-failure
RestartSec=2s

StandardOutput=journal
StandardError=journal

NoNewPrivileges=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictRealtime=yes

[Install]
WantedBy=multi-user.target
