# yocto/meta-iot/recipes-iot/lwm2m/files/iot-lwm2m-server.service
[Unit]
Description=IoT LwM2M Bootstrap + DM server
Documentation=https://github.com/naushada/iot
After=network-online.target iot-ds.service
Wants=network-online.target iot-ds.service
StartLimitBurst=5
StartLimitIntervalSec=60s

[Service]
Type=simple
EnvironmentFile=/etc/iot/lwm2m-server.env
ExecStart=/usr/bin/lwm2m $LWM2M_ARGS

DynamicUser=yes
# Join the `iot` group to open ds-server's 0660 group-iot socket.
SupplementaryGroups=iot

# No RuntimeDirectory=iot: /run/iot is owned by tmpfiles.d/iot.conf. Declaring
# it here would delete the shared dir (and ds-server's socket) on stop/restart.

Restart=on-failure
RestartSec=5s

StandardOutput=journal
StandardError=journal

NoNewPrivileges=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictRealtime=yes

[Install]
WantedBy=multi-user.target
