# yocto/meta-iot/recipes-iot/lwm2m/files/iot-wifi-client.service
[Unit]
Description=IoT wifi-client (wpa_supplicant + DHCP supervisor)
Documentation=https://github.com/naushada/iot
After=network-online.target iot-ds.service
Wants=network-online.target iot-ds.service
StartLimitBurst=5
StartLimitIntervalSec=60s

[Service]
Type=simple
EnvironmentFile=/etc/iot/wifi-client.env
ExecStart=/usr/bin/wifi-client $WIFI_CLIENT_ARGS

DynamicUser=yes
# Join the `iot` group to open ds-server's 0660 group-iot socket
# (/run/iot/data_store.sock).
SupplementaryGroups=iot
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW

DeviceAllow=/dev/rfkill rw

# No RuntimeDirectory=iot: /run/iot is owned by tmpfiles.d/iot.conf. Declaring
# it here would delete the shared dir (and ds-server's socket) on every
# stop/restart of this unit.

Restart=on-failure
RestartSec=5s

StandardOutput=journal
StandardError=journal

NoNewPrivileges=yes
PrivateTmp=no
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictRealtime=yes

[Install]
WantedBy=multi-user.target
