# yocto/meta-iot/recipes-iot/lwm2m/files/rpi3b-selftest.service
[Unit]
Description=rpi3B driver self-test (GoogleTest register bit-layout suite)
Documentation=https://github.com/naushada/rpi3B
After=local-fs.target
# Validation oneshot: runs the rpi3B GPIO/CLOCK/IRQ/I2C/SPI register-layout
# tests once at boot and records pass/fail in the journal. It exercises the
# bit math over heap-backed placement-new regions, so it needs no hardware
# access and does not poke real MMIO.

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/rpi3B_test --gtest_brief=1
StandardOutput=journal
StandardError=journal

# No privileges / no device access required — the suite is pure userspace.
DynamicUser=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateDevices=yes

[Install]
WantedBy=multi-user.target
